# .agents/launchd/com.spacory.agents.cycle.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- Rendered from this template by .agents/launchd/install.sh. Runs ONE product
     cycle (issue creation/refinement) on a slow cadence, separate from the fast
     implement→review→resolve loop, so ticket creation can't flood the queue. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.spacory.agents.cycle</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/.agents/dispatch.sh</string>
    <string>cycle</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>

  <!-- Once a day at 09:00 local. Adjust to taste. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/.agents/logs/cycle.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/.agents/logs/cycle.log</string>

  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>5</integer>
</dict>
</plist>
